# dist/ami/files/scylla-setup.service
[Unit]
Description=Scylla Setup
After=network.target

[Service]
Type=oneshot
ExecStart=/usr/lib/scylla/scylla-setup.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
