# dist/redhat/centos_dep/scylla-env-1.0/profile.d/scylla.sh
export PATH=$PATH:/opt/scylla/bin:/opt/scylla/sbin
