# dist/common/dep/scylla-env-1.0/profile.d/scylla.csh
set path = ($PATH /opt/scylladb/bin /opt/scylladb/sbin)

