# dist/common/dep/scylla-env-1.0/profile.d/scylla.sh
export PATH=$PATH:/opt/scylladb/bin:/opt/scylladb/sbin
