# dist/common/systemd/node-exporter.service
[Unit]
Description=Node Exporter

[Service]
Type=simple
User=scylla
Group=scylla
ExecStart=/usr/bin/node_exporter  --collector.interrupts

[Install]
WantedBy=multi-user.target
