# dist/common/systemd/scylla-housekeeping.service.in
[Unit]
Description=Scylla Housekeeping

[Service]
Type=simple
User=scylla
Group=scylla
ExecStart=/usr/lib/scylla/scylla-housekeeping -q version

[Install]
WantedBy=multi-user.target
