# dist/common/systemd/scylla-housekeeping.timer.in
[Unit]
Description=Run Scylla Housekeeping daily
After=scylla-server.service
BindsTo=scylla-server.service

[Timer]
OnUnitActiveSec=1d

[Install]
WantedBy=timers.target
