# dist/common/systemd/scylla-housekeeping.timer
[Unit]
Description=Run Scylla Housekeeping daily
After=scylla-server.service
BindsTo=scylla-server.service

[Timer]
# set OnActiveSec to 3 to safely avoid issues/1846
OnActiveSec=3
OnUnitActiveSec=1d

[Install]
WantedBy=timers.target
