# dist/common/systemd/scylla-housekeeping-restart.service
[Unit]
Description=Scylla Housekeeping restart mode
After=network.target

[Service]
Type=simple
User=scylla
Group=scylla
ExecStart=/usr/lib/scylla/scylla-housekeeping --uuid-file /var/lib/scylla-housekeeping/housekeeping.uuid -q --repo-files '/etc/yum.repos.d/scylla*.repo' -c /etc/scylla.d/housekeeping.cfg version --mode r

[Install]
WantedBy=multi-user.target
