# dist/common/systemd/scylla-housekeeping.service
[Unit]
Description=Scylla Housekeeping

[Service]
Type=simple
User=scylla
Group=scylla
ExecStart=/usr/lib/scylla/scylla-housekeeping -q -c /etc/scylla.d/housekeeping.cfg version

[Install]
WantedBy=multi-user.target
