# dist/common/systemd/scylla-housekeeping-restart.service
[Unit]
Description=Scylla Housekeeping restart mode
After=network.target

[Service]
Type=simple
User=scylla
Group=scylla
EnvironmentFile=/etc/sysconfig/scylla-housekeeping
ExecStart=/opt/scylladb/scripts/scylla-housekeeping --uuid-file $UUID_FILE -q -c $CONFIG_FILE --repo-files $REPO_FILES version --mode r

[Install]
WantedBy=multi-user.target
