# dist/common/systemd/scylla-housekeeping.service
[Unit]
Description=Scylla Housekeeping
After=network.target

[Service]
Type=simple
User=scylla
Group=scylla
ExecStart=/usr/lib/scylla/scylla-housekeeping --uuid-file /etc/scylla.d/housekeeping.uuid -q -c /etc/scylla.d/housekeeping.cfg version --mode d

[Install]
WantedBy=multi-user.target
