# dist/docker/etc/supervisord.conf.d/scylla-node-exporter.conf
[program:scylla-node-exporter]
command=/opt/scylladb/supervisor/scylla-node-exporter.sh
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
