# dist/docker/etc/supervisord.conf.d/scylla-housekeeping.conf
[program:scylla-housekeeping]
command=/scylla-housekeeping-service.sh
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
