# dist/docker/redhat/etc/supervisord.conf.d/scylla-jmx.conf
[program:scylla-jmx]
command=/scylla-jmx-service.sh
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
