# dist/common/systemd/scylla-server.service.d/dependencies.conf
[Unit]
After=local-fs.target network-online.target
Requires=local-fs.target network-online.target
