# dist/redhat/systemd/scylla-server.service
[Unit]
Description=Scylla Server
After=network.target libvirtd.service

[Service]
Type=simple
EnvironmentFile=/etc/sysconfig/scylla-server
ExecStartPre=/usr/lib/scylla/scylla_prepare
ExecStart=/usr/lib/scylla/scylla_run
ExecStopPost=/usr/lib/scylla/scylla_stop
KillMode=process
Restart=no

[Install]
WantedBy=multi-user.target
