# dist/redhat/systemd/scylla-io-setup.service
[Unit]
Description=Scylla IO Setup
After=network.target

[Service]
Type=oneshot
ExecStart=/usr/lib/scylla/scylla_io_setup
RemainAfterExit=yes
TimeoutStartSec=1800
